# files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.
mkdir -p /tmp/root
mkdir -p /tmp/root/.config
cp -r /root/micro /tmp/root/.config/
cp -r /root/.oh-my-zsh /tmp/root/
cp /root/.zshrc /tmp/root/

#mkdir -p /tmp/sing-box/rule_set
#mkdir -p /tmp/sing-box/config
#cp -r /root/sing-box /tmp/
#/etc/init.d/sing-box start

exit 0
